# etc/supervisor/conf.d/mysql.conf
[program:mysql_startup]
command=bash -c 'service mysql start'
autorestart=false
directory=/root
