# etc/supervisor/conf.d/node_supervisor.conf
[program:supervisor]
command=bash -c 'node server.js'
autorestart=false
directory=/usr/src/app/
