# etc/supervisor/conf.d/omysqlo.conf
[program:start_mysql]
command=bash -c '/usr/src/app/start_mysql.sh'
autorestart=false
directory=/usr/src/app/
